# html/crontab
# m h  dom mon dow   command
45 19 * * * ~/certmake/html/cron-html-update.sh
